# fedora/overlay/etc/systemd/system/systemd-tmpfiles-setup.service.d/override.conf
# This is necessary so we can create btrfs subvolumes even though / isn't a
# subvolume when using bootc.

[Service]
Environment=SYSTEMD_TMPFILES_FORCE_SUBVOL=1
