# fedora/overlay/etc/systemd/system/fstrim.service.d/override.conf
# https://github.com/fedora-silverblue/issue-tracker/issues/689

[Service]
ExecStart=/usr/bin/fstrim --listed-in /proc/self/mountinfo --verbose --quiet-unsupported
